# _units/sample.service
[Unit]

[Service]
ExecStart=echo "Hello, World!"
ExecStop=echo "Goodbye, World!"

[Install]
